# windfoil-proxy (1).service
# windfoil-proxy.service — Node-Proxy mit Admin-Funktionen
# Install: nach /etc/systemd/system/ kopieren, dann:
#   sudo systemctl daemon-reload && sudo systemctl enable --now windfoil-proxy
#
# Secrets stehen in /var/www/windfoil/windfoil.env (NICHT hier, NICHT in Git).

[Unit]
Description=WindFoil Weather Station Proxy
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/windfoil
EnvironmentFile=/var/www/windfoil/windfoil.env
ExecStart=/usr/bin/node proxy-server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
